# hw04/ansible/nodebuntu/roles/setup-hardening-for-cis-benchmarks/tasks/autz-broker/twistlock-authz.service

[Unit]
Description=Twistlock docker authorization plugin
After=syslog.target
[Service]
Type=simple
ExecStart=/usr/bin/authz-broker
[Install]
WantedBy=multi-user.target